# foodmaps.service
[Unit]
Description=Food Maps API Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/project/backend

# Environment
Environment="PATH=/home/ec2-user/.local/bin:/home/ec2-user/project/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ec2-user/project"
EnvironmentFile=/home/ec2-user/project/.env

# Start script
ExecStart=/bin/bash /home/ec2-user/project/backend/start_server.sh

# Auto-restart configuration
Restart=always
RestartSec=10
# Bound restart attempts to avoid thundering-herd / CPU spin if the
# process crashes immediately on boot (e.g. config error). Allow up to
# 5 restarts per 60s before giving up.
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foodmaps

[Install]
WantedBy=multi-user.target
